# recipes-atlas-soc/atlas-soc-fftsw-apps/files/atlas-soc-fftsw-init.service
[Unit]
Description=Atlas SoC Prep fftsw

[Service]
ExecStart=/examples/fft/bin/setup_target_fft_env.sh

[Install]
WantedBy=multi-user.target
